# nmeamux/root/etc/cron.daily/nmemamux-restart
#!/bin/sh

ps | grep 'lua.agent.*nmeamux' | grep -v grep | awk '{print $1}' | xargs kill
